# validation/hardware/hats/filesystem/etc/profile.d/HatsProfile.sh
# HATS environment variable

export LD_LIBRARY_PATH=/lib:/usr/lib:/system/lib:/system/usr/lib:/system/lib/ste_omxcomponents
export PATH=/bin:/sbin:/usr/bin:/usr/sbin:/system/bin:/system/usr/bin
