# validation/hardware/hats/filesystem/etc/profile.d/mmteprofile.sh
# MMTE environment variable
export OMX_BELLAGIO_REGISTRY=/system/omxcomponents
